# droppppn/supervisor/conf.d/php-fpm-71.conf
[program:php-fpm-71]
command=/opt/app/php71/sbin/php-fpm --nodaemonize --fpm-config /etc/php71/php-fpm.conf
user=root
autostart=true
priority=98
startretries=3
autorestart=true
stdout_logfile=/data/logs/supervisor/%(program_name)s.log
stderr_logfile=/data/logs/supervisor/%(program_name)s_err.log
