# dropppmn/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/bin/pidproxy /data/mysql/mysql.pid /usr/bin/mysqld_safe
user=mysql
autostart=true
priority=2
startretries=3
autorestart=true
stdout_logfile=/data/logs/supervisor/%(program_name)s.log
stderr_logfile=/data/logs/supervisor/%(program_name)s_err.log
