# droppn/supervisor/conf.d/mysql.conf
[program:mysqld_safe]
command=/usr/bin/mysqld_safe
process_name=%(program_name)s
numprocs=1
user=mysql
autostart=true
priority=97
startretries=9999
autorestart=true
stdout_logfile=/data/logs/supervisor/%(program_name)s.log
stderr_logfile=/data/logs/supervisor/%(program_name)s_err.log